# rootfs/etc/services.d/node-inspector/run
#!/usr/bin/execlineb -P

with-contenv
backtick -D "" -n START_NODE_INSPECTOR_ON_PORT { printcontenv START_NODE_INSPECTOR_ON_PORT }
importas -u START_NODE_INSPECTOR_ON_PORT START_NODE_INSPECTOR_ON_PORT
ifelse { s6-test -n ${START_NODE_INSPECTOR_ON_PORT} }
{
	fdmove -c 2 1
  node-inspector --web-port=${START_NODE_INSPECTOR_ON_PORT}
}
